# templates/btrfs-scrub.timer.j2
[Unit]
Description=Scrub BTRFS disks every month

[Timer]
Persistent=yes
OnCalendar={{ btrfs_scrub_oncalendar }}

[Install]
WantedBy=timers.target
