# files/btrfs-scrub.service
[Unit]
Description=Scrub all BTRFS disks

[Service]
Type=simple
ExecStart=/usr/local/bin/btrfs-scrub
ExecStop=-/usr/local/bin/btrfs-scrub-cancel
